# assembler/payload/unity-net.service
[Unit]
Description=Unity network configuration service

[Service]
ExecStart=/root/unity-net.sh

[Install]
WantedBy=multi-user.target
