# transducers/processkiller/etc/systemd/system/processkiller.service
[Unit]
After=network.target

[Service]
Type=forking
ExecStart=/etc/init.d/processkiller start

[Install]
WantedBy=multi-user.target
